# software/irrigatore-daemon.service
# Demone di irrigazione: unico proprietario delle pompe/GPIO.
# Copia in /etc/systemd/system/ e adatta i percorsi (User, WorkingDirectory, ExecStart).
#   sudo cp software/irrigatore-daemon.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now irrigatore-daemon
#   journalctl -u irrigatore-daemon -f
[Unit]
Description=Irrigazione balcone - demone pompe
After=network.target

[Service]
Type=simple
User=mzaccari
WorkingDirectory=/home/mzaccari/irrigatore
Environment=PUMP_MOCK=0
Environment=SENSOR_MOCK=0
Environment=GPIOZERO_PIN_FACTORY=lgpio
# Il topic ntfy (segreto) NON va qui: usa un drop-in non versionato, vedi
# docs/upgrade-sensori.md ("notify.conf" con NOTIFY_NTFY_TOPIC).
ExecStart=/home/mzaccari/irrigatore/.venv/bin/python /home/mzaccari/irrigatore/software/daemon.py
# All'avvio e alla chiusura il demone porta tutte le pompe in stato sicuro (spente).
Restart=always
RestartSec=2
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
